# rootfs/etc/s6-overlay/s6-rc.d/searxng-configure/up
#!/usr/bin/env sh

: ${CATEGORIES}
: ${DEFAULT_ENGINES}
: ${ENABLE_DONATION_URL}
: ${ENABLE_METRICS}
: ${IS_PUBLIC_INSTANCE}
: ${SECRET_KEY}

cp /usr/local/searxng/searx/limiter.toml /etc/searxng/limiter.toml

yq -i ".general.enable_metrics = env(ENABLE_METRICS)"     /etc/searxng/settings.yml
yq -i ".plugins[].active = false"                         /usr/local/searxng/searx/settings.yml
yq -i ".server.public_instance = env(IS_PUBLIC_INSTANCE)" /etc/searxng/settings.yml
yq -i ".server.secret_key = env(SECRET_KEY)"              /etc/searxng/settings.yml

unset SECRET_KEY

if [ -n "$CATEGORIES" ]; then
    yq -i '.categories_as_tabs = {}
        | (
            strenv(CATEGORIES)
            | split(",")
            | map(trim)
        ) as $x
        | .categories_as_tabs[$x[]] = ~
    ' /etc/searxng/settings.yml
else
    yq -i 'del(.categories_as_tabs)' /etc/searxng/settings.yml
fi

if [ -n "$DEFAULT_ENGINES" ]; then
    yq -i ".engines[].disabled = true" /usr/local/searxng/searx/settings.yml

    yq -i '.engines += (
        strenv(DEFAULT_ENGINES)
        | split(",")
        | map(trim)
        | map({
            "name": .,
            "disabled": false
        })
    )' /etc/searxng/settings.yml
else
    yq -i 'del(.engines)' /etc/searxng/settings.yml
fi
