# rootfs/etc/s6-overlay/s6-rc.d/valkey-server/run
#!/usr/bin/env sh

exec /usr/bin/valkey-server \
    /etc/valkey/valkey.conf
